# root/etc/s6-overlay/s6-rc.d/init-agents/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "**** agents-manager init: installing claude agents ****"

# create aatos user with passwordless sudo (if not exists)
if ! id "aatos" &>/dev/null; then
    useradd -m -u 911 -g abc -s /bin/bash aatos
    echo "aatos ALL=(ALL) NOPASSWD: ALL" > /etc/sudoers.d/aatos
    chmod 440 /etc/sudoers.d/aatos
    echo "root ALL=(ALL) NOPASSWD: ALL" > /etc/sudoers.d/root
    chmod 440 /etc/sudoers.d/root
    echo "created aatos account"
fi

# Install claude for root if not present
if ! command -v claude &>/dev/null; then
    echo "**** installing claude for root ****"
    curl -fsSL https://claude.ai/install.sh | bash
else
    echo "claude already installed for root"
fi

# write root .claude/settings.json using python
python3 << 'PYEOF'
import os, json
settings = {
    "skipDangerousModePermissionPrompt": True,
    "env": {
        "ANTHROPIC_BASE_URL": os.environ.get("MINIMAX_ANTHROPIC_BASE_URL", ""),
        "ANTHROPIC_MODEL": os.environ.get("LLM_MODEL", ""),
        "ANTHROPIC_DEFAULT_SONNET_MODEL": os.environ.get("LLM_MODEL", ""),
        "ANTHROPIC_DEFAULT_OPUS_MODEL": os.environ.get("LLM_MODEL", ""),
        "ANTHROPIC_DEFAULT_HAIKU_MODEL": os.environ.get("LLM_MODEL", ""),
        "CLAUDE_CODE_SUBAGENT_MODEL": os.environ.get("LLM_MODEL", ""),
        "ANTHROPIC_AUTH_TOKEN": os.environ.get("ANTHROPIC_API_KEY", ""),
        "CLAUDE_CODE_EXPERIMENTAL_AGENT_TEAMS": "10",
        "teammateMode": "tmux",
        "MINIMAX_API_KEY": os.environ.get("MINIMAX_API_KEY", ""),
        "GITHUB_TOKEN": os.environ.get("GITHUB_TOKEN", ""),
        "MATTERMOST_URL": os.environ.get("MATTERMOST_URL", ""),
        "MATTERMOST_TOKEN": os.environ.get("MATTERMOST_TOKEN", "")
    },
    "dangerouslyAlwaysAllow": True,
    "allow": ["Edit", "Write", "Bash", "Read", "Glob", "Grep", "WebFetch", "WebSearch", "TodoRead", "TodoWrite"]
}
os.makedirs("/root/.claude/settings", exist_ok=True)
with open("/root/.claude/settings.json", "w") as f:
    json.dump(settings, f, indent=2)
print("wrote /root/.claude/settings.json")
PYEOF

# add PATH to root bashrc
if ! grep -q '.local/bin' /root/.bashrc 2>/dev/null; then
    echo 'export PATH="$HOME/.local/bin:$PATH"' >> /root/.bashrc
fi

# create .claude.json for root (skip onboarding)
python3 -c 'import json; json.dump({"hasCompletedOnboarding": True}, open("/root/.claude.json", "w"))'

# install claude for aatos if not present
if ! su - aatos -c "command -v claude" &>/dev/null; then
    echo "**** installing claude for aatos ****"
    su - aatos -c "curl -fsSL https://claude.ai/install.sh | bash"
else
    echo "claude already installed for aatos"
fi

# write aatos .claude/settings.json using python
su - aatos -c "python3" << 'PYEOF'
import os, json
settings = {
    "skipDangerousModePermissionPrompt": True,
    "env": {
        "ANTHROPIC_BASE_URL": os.environ.get("MINIMAX_ANTHROPIC_BASE_URL", ""),
        "ANTHROPIC_MODEL": os.environ.get("LLM_MODEL", ""),
        "ANTHROPIC_DEFAULT_SONNET_MODEL": os.environ.get("LLM_MODEL", ""),
        "ANTHROPIC_DEFAULT_OPUS_MODEL": os.environ.get("LLM_MODEL", ""),
        "ANTHROPIC_DEFAULT_HAIKU_MODEL": os.environ.get("LLM_MODEL", ""),
        "CLAUDE_CODE_SUBAGENT_MODEL": os.environ.get("LLM_MODEL", ""),
        "ANTHROPIC_AUTH_TOKEN": os.environ.get("ANTHROPIC_API_KEY", ""),
        "CLAUDE_CODE_EXPERIMENTAL_AGENT_TEAMS": "10",
        "teammateMode": "tmux",
        "MINIMAX_API_KEY": os.environ.get("MINIMAX_API_KEY", ""),
        "GITHUB_TOKEN": os.environ.get("GITHUB_TOKEN", ""),
        "MATTERMOST_URL": os.environ.get("MATTERMOST_URL", ""),
        "MATTERMOST_TOKEN": os.environ.get("MATTERMOST_TOKEN", "")
    },
    "dangerouslyAlwaysAllow": True,
    "allow": ["Edit", "Write", "Bash", "Read", "Glob", "Grep", "WebFetch", "WebSearch", "TodoRead", "TodoWrite"]
}
os.makedirs(os.path.expanduser("~/.claude/settings"), exist_ok=True)
with open(os.path.expanduser("~/.claude/settings.json"), "w") as f:
    json.dump(settings, f, indent=2)
print("wrote ~/.claude/settings.json")
PYEOF

su - aatos -c 'echo "export PATH=\"$HOME/.local/bin:$PATH\"" >> ~/.bashrc'
su - aatos -c 'echo "{\"hasCompletedOnboarding\": true}" > ~/.claude.json'

# create agents-manager launchers in /usr/local/bin
if ln -sf /config/agents-manager/actions.sh /usr/local/bin/agents-manager 2>/dev/null; then
    echo "created /usr/local/bin/agents-manager"
fi
if ln -sf /config/agents-manager/actions.sh /usr/local/bin/actions 2>/dev/null; then
    echo "created /usr/local/bin/actions"
fi

# symlink workspace to aatos home
if [ -d "/config/workspace" ] && [ ! -L "/home/aatos" ]; then
    ln -sf /config/workspace /home/aatos 2>/dev/null || true
fi

echo "**** agents-manager init complete ****"
